# units/unifi-controller@.service.d/publish-8443.conf
[Service]
# Expose the HTTPS port (8443)
Environment='DOCKER_PUBLISH_8443=8443:8443'
