# units/unifi-controller@.service.d/publish-8080.conf
[Service]
# Expose the HTTP port (8080)
Environment='DOCKER_PUBLISH_8080=8080:8080'
